# zeitkonto.ini
[program:zeitkonto]
directory=/var/www/virtual/drevil/zeitkonto/workforce_management
command=/home/drevil/.local/bin/gunicorn --error-logfile - --reload --bind 0.0.0.0:8000 workforce_management.wsgi:application
startsecs=15
    
